# ai-necklace-openai.service
[Unit]
Description=AI Necklace OpenAI Realtime Voice Assistant v5 (Autonomous Mode)
After=network-online.target sound.target
Requires=network-online.target

[Service]
Type=simple
User=hotaka
WorkingDirectory=/home/hotaka/raspi-voice5
Environment=PATH=/home/hotaka/raspi-voice5/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/hotaka/.ai-necklace/.env
ExecStart=/home/hotaka/raspi-voice5/venv/bin/python3 ai_necklace_openai.py

# Auto-restart configuration
Restart=always
RestartSec=10
StartLimitIntervalSec=600
StartLimitBurst=20

# Graceful shutdown
TimeoutStopSec=15

# Logging
StandardOutput=journal
StandardError=journal

# Audio + GPIO access
SupplementaryGroups=audio gpio

[Install]
WantedBy=multi-user.target
